# deploy/raspberry-pi/mrrfantasy-frontend.service
[Unit]
Description=MRR Fantasy Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/MrrFantasy/frontend
EnvironmentFile=/etc/mrrfantasy-frontend.env
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
